# systemd/.config/systemd/user/box.service
[Unit]
Description=File Sync Service

[Service]
Type=oneshot
ExecStart=/usr/bin/unison box -batch

[Install]
WantedBy=network.target
