# systemd/.config/systemd/user/dav.service
[Unit]
Description=DAV Sync Service

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync

[Install]
WantedBy=network.target
