# systemd/.config/systemd/user/feed.service
[Unit]
Description=Feed Service

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -X reload

[Install]
WantedBy=network.target
